# garlic/bus/1.0/config/supervisor.conf
[program:garlic.bus]
command=/var/www/domain.com/bin/console --env=prod --no-debug enqueue:consume --setup-broker
process_name=%(program_name)s_%(process_num)02d
numprocs=4
autostart=true
autorestart=true
tartsecs=0
user=www-data
redirect_stderr=true